# daemons/eskil-telemetry.service
[Unit]
Description=ARES Eskil Telemetry Bridge Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /root/daemons/arduino_uno_q_telemetry.py
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
